# IC/Digital/Verification/sv/7/run
#!/bin/sh
irun +access+wrc test.sv -input wave.tcl
